# files/bypass-gateway/etc/uci-defaults/99-bypass-gateway
#!/bin/sh
# Bypass-gateway role first-boot:
#   - Static LAN, no DHCP server (main router serves DHCP)
#   - No WAN
#   - dnsmasq → sing-box DNS (anti-leak, sing-box's own routing)
#   - homeproxy default to bypass-mainland-china mode

# --- LAN: static, default to .10.2 (user changes in LuCI as needed) ---
uci -q batch <<EOF
set network.lan.proto='static'
set network.lan.ipaddr='192.168.10.2'
set network.lan.netmask='255.255.255.0'
set network.lan.gateway='192.168.10.1'
set network.lan.dns='192.168.10.1'
commit network
EOF

# --- Remove WAN (bypass has no WAN) ---
uci -q delete network.wan
uci -q delete network.wan6
uci commit network

# --- DHCP server OFF (main router serves DHCP) ---
uci -q set dhcp.lan.ignore='1'
uci commit dhcp

# --- DNS upstream → sing-box (zero-leak split routing) ---
uci -q delete dhcp.@dnsmasq[0].server
uci -q add_list dhcp.@dnsmasq[0].server='127.0.0.1#5333'
uci commit dhcp

# --- homeproxy default settings ---
[ -f /etc/config/homeproxy ] && uci -q batch <<EOF
set homeproxy.config.routing_mode='bypass_mainland_china'
set homeproxy.config.ipv6_support='0'
set homeproxy.config.china_dns_server='223.5.5.5'
commit homeproxy
EOF

exit 0
