# files/common/etc/uci-defaults/91-cert-gen
#!/bin/sh
# Generate self-signed CA + server cert for nginx (LuCI HTTPS).
# Designed for internal LAN use where ACME is impractical (CN ISPs block 80/443).
# Cert is valid 10 years; covers IP, hostname, and *.lan.
# CA cert is exposed at http://<router-ip>/ca.crt for client trust import.

set -e

CERT_DIR=/etc/nginx/conf.d
CA_DIR=/etc/ssl/iwrt-ca
WWW_DIR=/www

mkdir -p "$CERT_DIR" "$CA_DIR" "$WWW_DIR"

# Skip if already generated (handles re-run / sysupgrade)
[ -f "$CERT_DIR/_lan.crt" ] && [ -f "$WWW_DIR/ca.crt" ] && exit 0

# --- Pull dynamic identity from current config ---
LAN_IP=$(uci -q get network.lan.ipaddr || echo "192.168.10.1")
HOSTNAME=$(uci -q get system.@system[0].hostname || echo "ImmortalWrt-Router")

# --- Generate ECDSA P-256 CA (key + self-signed root) ---
openssl ecparam -name prime256v1 -genkey -noout -out "$CA_DIR/ca.key"
openssl req -new -x509 -days 3650 -key "$CA_DIR/ca.key" -out "$CA_DIR/ca.crt" \
    -subj "/CN=ImmortalWrt Internal CA/O=$HOSTNAME"

# --- Generate server key + CSR + sign with CA ---
openssl ecparam -name prime256v1 -genkey -noout -out "$CERT_DIR/_lan.key"

cat > "$CA_DIR/server.cnf" <<EOF
[req]
distinguished_name = req_dn
prompt = no
[req_dn]
CN = ${HOSTNAME}
O  = ${HOSTNAME}
[v3_ext]
subjectAltName = @alt
keyUsage = critical, digitalSignature, keyEncipherment
extendedKeyUsage = serverAuth
basicConstraints = CA:FALSE
[alt]
DNS.1 = ${HOSTNAME}
DNS.2 = router.lan
DNS.3 = gw.lan
DNS.4 = *.lan
IP.1  = ${LAN_IP}
IP.2  = 127.0.0.1
EOF

openssl req -new -key "$CERT_DIR/_lan.key" -out "$CA_DIR/server.csr" -config "$CA_DIR/server.cnf"
openssl x509 -req -in "$CA_DIR/server.csr" \
    -CA "$CA_DIR/ca.crt" -CAkey "$CA_DIR/ca.key" -CAcreateserial \
    -out "$CERT_DIR/_lan.crt" -days 3650 \
    -extfile "$CA_DIR/server.cnf" -extensions v3_ext

# --- Permissions ---
chmod 600 "$CERT_DIR/_lan.key" "$CA_DIR/ca.key"
chmod 644 "$CERT_DIR/_lan.crt" "$CA_DIR/ca.crt"

# --- Expose CA over HTTP for client trust download ---
cp "$CA_DIR/ca.crt" "$WWW_DIR/ca.crt"
chmod 644 "$WWW_DIR/ca.crt"

# --- nginx HTTP server: allow plain-HTTP /ca.crt download, redirect rest -> HTTPS
mkdir -p /etc/nginx/conf.d
cat > /etc/nginx/conf.d/_ca-cert.conf <<'NGX'
# Auto-generated by 91-cert-gen.
# Allows clients to download the CA over plain HTTP (so they can trust the
# self-signed cert before browsing HTTPS LuCI). Everything else -> HTTPS.
server {
    listen 80 default_server;
    listen [::]:80 default_server;
    location = /ca.crt {
        root /www;
        default_type application/x-x509-ca-cert;
        add_header Content-Disposition 'attachment; filename="iwrt-ca.crt"';
    }
    location / {
        return 301 https://$host$request_uri;
    }
}
NGX

# --- Compute fingerprint and stamp into banner ---
FPR=$(openssl x509 -noout -fingerprint -sha256 -in "$CA_DIR/ca.crt" | sed 's/SHA256 Fingerprint=//')
echo "$FPR" > "$CA_DIR/ca.fingerprint"
sed -i "s|%CA_FPR%|$FPR|g" /etc/banner 2>/dev/null || true

# --- Drop CA cleanup ---
rm -f "$CA_DIR/server.csr" "$CA_DIR/server.cnf"

logger -t cert-gen "self-signed cert generated for $HOSTNAME ($LAN_IP); CA fingerprint $FPR"
exit 0
