# files/common/etc/uci-defaults/92-netdata-harden
#!/bin/sh
# Harden netdata for LAN access:
#   - Bind to LAN IP + 127.0.0.1 (not all interfaces)
#   - Allow web access only from RFC1918 LAN ranges
#   - Disable cloud, registry, anonymous stats
#   - Firewall REJECT :19999 from WAN zone (LuCI iframe still works on LAN)
#
# Trade-off: netdata 1.38 has no built-in auth, so anyone on LAN can view.
# For an enterprise edge router on a trusted LAN that's acceptable; if you
# want stricter, use 'ssh -L 19999:localhost:19999 root@router' tunnel.

set -e

NETDATA_CONF=/etc/netdata/netdata.conf

mkdir -p /etc/netdata
cat > "$NETDATA_CONF" <<'EOF'
[global]
    run as user = netdata
    memory mode = ram
    history = 3600
    update every = 2
    anonymous statistics = no

[web]
    bind to = *
    allow connections from = localhost 192.168.*.* 10.*.*.* 172.16.*.* 172.17.*.* 172.18.*.* 172.19.*.* 172.20.*.* 172.21.*.* 172.22.*.* 172.23.*.* 172.24.*.* 172.25.*.* 172.26.*.* 172.27.*.* 172.28.*.* 172.29.*.* 172.30.*.* 172.31.*.*
    allow dashboard from = localhost 192.168.*.* 10.*.*.* 172.16.*.* 172.17.*.* 172.18.*.* 172.19.*.* 172.20.*.* 172.21.*.* 172.22.*.* 172.23.*.* 172.24.*.* 172.25.*.* 172.26.*.* 172.27.*.* 172.28.*.* 172.29.*.* 172.30.*.* 172.31.*.*
    allow badges from = *
    allow streaming from = none
    allow netdata.conf from = localhost
    allow management from = none
    enable gzip compression = yes

[registry]
    enabled = no

[cloud]
    enabled = no

[plugins]
    cgroups = no
    apps = no
    proc = yes
    diskspace = yes
    tc = no
    idlejitter = no
    python.d = no
    charts.d = no
    node.d = no
    go.d = no
EOF

# Firewall: REJECT :19999 from WAN zone only (NOT '*' — that would block LAN too,
# breaking LuCI's Netdata iframe). LAN access remains.
uci -q batch <<EOF
add firewall rule
set firewall.@rule[-1].name='Block-Netdata-WAN'
set firewall.@rule[-1].src='wan'
set firewall.@rule[-1].proto='tcp'
set firewall.@rule[-1].dest_port='19999'
set firewall.@rule[-1].target='REJECT'
commit firewall
EOF

/etc/init.d/netdata enable 2>/dev/null || true

logger -t netdata-harden "netdata bound to all interfaces; firewall REJECTs WAN -> :19999"
exit 0
