# files/common/etc/uci-defaults/92-netdata-harden
#!/bin/sh
# Harden netdata: bind to localhost only, disable cloud/registry/anonymous
# stats. External access is via SSH tunnel (ssh -L 19999:localhost:19999
# root@router) then http://localhost:19999/ in a browser.
#
# We deliberately do NOT set up a reverse proxy. The prior attempt with
# nginx was over-engineered and broken (missing packages in 24.10 feed).
# LAN-bound + firewall REJECT + SSH-tunnel-only is simpler and auditable.

set -e

NETDATA_CONF=/etc/netdata/netdata.conf

mkdir -p /etc/netdata
cat > "$NETDATA_CONF" <<'EOF'
[global]
    run as user = netdata
    memory mode = ram
    history = 3600
    update every = 2
    anonymous statistics = no

[web]
    bind to = 127.0.0.1
    allow connections from = localhost
    allow dashboard from = localhost
    allow badges from = *
    allow streaming from = none
    allow netdata.conf from = localhost
    allow management from = none
    enable gzip compression = yes

[registry]
    enabled = no

[cloud]
    enabled = no

[plugins]
    cgroups = no
    apps = no
    proc = yes
    diskspace = yes
    tc = no
    idlejitter = no
    python.d = no
    charts.d = no
    node.d = no
    go.d = no
EOF

# Firewall: ensure :19999 can't leak even if bind is ever widened later
uci -q batch <<EOF
add firewall rule
set firewall.@rule[-1].name='Block-Netdata-External'
set firewall.@rule[-1].src='*'
set firewall.@rule[-1].proto='tcp'
set firewall.@rule[-1].dest_port='19999'
set firewall.@rule[-1].target='REJECT'
commit firewall
EOF

/etc/init.d/netdata enable 2>/dev/null || true

logger -t netdata-harden "netdata bound to 127.0.0.1. Access via 'ssh -L 19999:localhost:19999 root@router' then http://localhost:19999/"
exit 0
