# files/common/etc/uci-defaults/93-services-disable
#!/bin/sh
# Services that require user-supplied configuration AND would actively
# misbehave or leak when started unconfigured. Disabled at first boot.
#
# IMPORTANT: client-side proxy frameworks (passwall / homeproxy /
# sing-box / xray) are NOT disabled. They start with empty config and
# idle harmlessly until the user imports a node via LuCI. Disabling
# their init link breaks the LuCI "enable" toggle (which only flips UCI,
# does not invoke init.d start), making the user think the firmware is
# broken. Hard lesson learned from .202 deployment.

set -e

# passwall_server is a SERVER (would expose vmess/vless port to LAN).
# Always off until the user explicitly configures and starts it.
SERVER_SERVICES="passwall_server"

# Discovery / file-sharing helpers that broadcast router presence on LAN
DISCOVERY_SERVICES="wsdd2"

# Tunnel / DDNS / web-shell — useless without external config
TUNNEL_SERVICES="frpc mwan3 ttyd freeradius3 strongswan"

# Pulled in as a passwall dependency, never used standalone
UNUSED_SERVICES="haproxy"

ALL="$SERVER_SERVICES $DISCOVERY_SERVICES $TUNNEL_SERVICES $UNUSED_SERVICES"

DISABLED=""
for svc in $ALL; do
    if [ -x "/etc/init.d/$svc" ]; then
        "/etc/init.d/$svc" disable 2>/dev/null || true
        "/etc/init.d/$svc" stop    2>/dev/null || true
        DISABLED="$DISABLED $svc"
    fi
done

logger -t services-disable "init disabled at boot:${DISABLED:- (none installed)}"
exit 0
