# files/common/etc/uci-defaults/99-common-hardening
#!/bin/sh
# Runs once on first boot. Hardening common to all roles
# (main-router / bypass-gateway / mtk-filogic-enterprise).

# --- dropbear: key-only, no root password login ---
uci -q batch <<EOF
set dropbear.@dropbear[0].PasswordAuth='off'
set dropbear.@dropbear[0].RootPasswordAuth='off'
set dropbear.@dropbear[0].GatewayPorts='off'
commit dropbear
EOF

# --- dnsmasq hardening: forward to local https-dns-proxy (DoH) ---
# Upstream DoH listeners are configured by 92-netdata-harden's siblings
# (https-dns-proxy ships defaulting to 127.0.0.1:5053 alidns DoH +
#  127.0.0.1:5054 cloudflare DoH).
uci -q batch <<EOF
set dhcp.@dnsmasq[0].noresolv='1'
set dhcp.@dnsmasq[0].rebind_protection='1'
set dhcp.@dnsmasq[0].rebind_localhost='1'
set dhcp.@dnsmasq[0].localservice='1'
set dhcp.@dnsmasq[0].cachesize='2000'
set dhcp.@dnsmasq[0].filter_aaaa='1'
set dhcp.@dnsmasq[0].interface='br-lan'
set dhcp.@dnsmasq[0].domainneeded='1'
delete dhcp.@dnsmasq[0].server
add_list dhcp.@dnsmasq[0].server='127.0.0.1#5053'
add_list dhcp.@dnsmasq[0].server='127.0.0.1#5054'
commit dhcp
EOF

# --- https-dns-proxy: ensure DoH instances are configured ---
if [ -x /etc/init.d/https-dns-proxy ]; then
  uci -q delete https-dns-proxy 2>/dev/null
  uci -q batch <<EOF
add https-dns-proxy https-dns-proxy
set https-dns-proxy.@https-dns-proxy[-1].bootstrap_dns='223.5.5.5'
set https-dns-proxy.@https-dns-proxy[-1].resolver_url='https://dns.alidns.com/dns-query'
set https-dns-proxy.@https-dns-proxy[-1].listen_addr='127.0.0.1'
set https-dns-proxy.@https-dns-proxy[-1].listen_port='5053'
add https-dns-proxy https-dns-proxy
set https-dns-proxy.@https-dns-proxy[-1].bootstrap_dns='1.1.1.1'
set https-dns-proxy.@https-dns-proxy[-1].resolver_url='https://1.1.1.1/dns-query'
set https-dns-proxy.@https-dns-proxy[-1].listen_addr='127.0.0.1'
set https-dns-proxy.@https-dns-proxy[-1].listen_port='5054'
commit https-dns-proxy
EOF
fi

# --- Firewall: disable hardware/software fastpath + fullcone ---
# These bypass nftables chain once a flow is established, breaking
# transparent proxy frameworks (passwall / openclash / sing-box tproxy).
# Verified live on RAX3000M: with fastpath on, PSW_DNS counter stayed
# at 0 packets while client browsing failed.
uci -q batch <<EOF
set firewall.@defaults[0].flow_offloading='0'
set firewall.@defaults[0].flow_offloading_hw='0'
set firewall.@defaults[0].fullcone='0'
set firewall.@defaults[0].fullcone6='0'
set firewall.@defaults[0].synflood_protect='1'
commit firewall
EOF

# --- Disable noisy / risky services at boot ---
# (also covered by 93-services-disable; redundant but harmless)
for s in ttyd upnpd vlmcsd spotifyd rustdesk-server uhttpd; do
  [ -x /etc/init.d/$s ] && /etc/init.d/$s disable 2>/dev/null
done

# --- nginx becomes the sole web server (uhttpd disabled above) ---
# luci-ssl-nginx gives us HTTPS LuCI; uhttpd stays installed but inert.
# Self-signed cert is generated by 91-cert-gen, applied to nginx's
# /etc/nginx/conf.d/_lan.conf.

# --- IPv6: explicitly disable RA / DHCPv6 on LAN by default ---
uci -q batch <<EOF
set network.lan.ipv6='0'
set dhcp.lan.dhcpv6='disabled'
set dhcp.lan.ra='disabled'
set dhcp.lan.ndp='disabled'
commit network
commit dhcp
EOF

# --- Banner: identify build ---
cat > /etc/banner <<'BANNER'
  ___                            _   _ __        __    _
 |_ _|_ __ ___  _ __ ___   ___  | |_| |\ \      / / __| |_
  | || '_ ` _ \| '_ ` _ \ / _ \ | __| __\ \ /\ / / '__| __|
  | || | | | | | | | | | | (_) || |_| |_ \ V  V /| |  | |_
 |___|_| |_| |_|_| |_| |_|\___/  \__|\__| \_/\_/ |_|   \__|
   Production build (hardened) — built %BUILDDATE%

   * LuCI: https://<router-ip>/   (uhttpd disabled, nginx HTTPS only)
   * CA:   http://<router-ip>/ca.crt   (import to trust the cert)
   * netdata: https://<router-ip>/netdata/   (creds in /root/.netdata-password)
   * SSH: key-only by default. Add key via LuCI before next reboot.
BANNER

exit 0
