# files/main-router/etc/uci-defaults/99-main-router
#!/bin/sh
# Main-router role first-boot:
#   - LAN 192.168.10.1, DHCP server on
#   - WAN scaffold (eth1, DHCP — user changes to PPPoE in LuCI)
#   - WAN firewall hardened
#   - dnsmasq → sing-box DNS (anti-leak)
#   - UPnP installed but disabled by default

# --- LAN ---
uci -q batch <<EOF
set network.lan.ipaddr='192.168.10.1'
set network.lan.netmask='255.255.255.0'
commit network
EOF

# --- WAN scaffold (DHCP placeholder — user picks PPPoE in LuCI) ---
uci -q batch <<EOF
delete network.wan
set network.wan=interface
set network.wan.proto='dhcp'
set network.wan.device='eth1'
set network.wan.peerdns='0'
set network.wan.dns='223.5.5.5 119.29.29.29'
commit network
EOF

# --- Firewall: WAN zone default REJECT (commercial standard) ---
uci -q batch <<EOF
set firewall.@zone[1].input='REJECT'
set firewall.@zone[1].output='ACCEPT'
set firewall.@zone[1].forward='REJECT'
set firewall.@zone[1].masq='1'
set firewall.@zone[1].mtu_fix='1'
commit firewall
EOF

# --- DNS upstream → sing-box (zero-leak split routing when proxy active) ---
uci -q delete dhcp.@dnsmasq[0].server
uci -q add_list dhcp.@dnsmasq[0].server='127.0.0.1#5333'
uci commit dhcp

# --- UPnP installed but DISABLED (enable in LuCI when needed) ---
[ -x /etc/init.d/miniupnpd ] && /etc/init.d/miniupnpd disable 2>/dev/null

exit 0
