# files/mtk-filogic-enterprise/etc/uci-defaults/92-netdata-harden
#!/bin/sh
# Harden netdata: bind LAN-only, disable cloud/registry, generate Basic Auth password,
# and front it with nginx reverse proxy so external access requires HTTPS+auth.

set -e

NETDATA_CONF=/etc/netdata/netdata.conf
HTPASSWD=/etc/nginx/.netdata.htpasswd
PASS_FILE=/root/.netdata-password
NGINX_LOC=/etc/nginx/conf.d/netdata.locations

# --- Generate strong random password ---
PASS=$(head -c 18 /dev/urandom | base64 | tr -d '+/=' | head -c 24)
echo "admin:$(openssl passwd -apr1 "$PASS")" > "$HTPASSWD"
chmod 600 "$HTPASSWD"

# Save password for admin to read once
{
    echo "Netdata access:"
    echo "  URL:      https://router.lan/netdata/"
    echo "  Username: admin"
    echo "  Password: $PASS"
    echo ""
    echo "(Generated $(date -u +%Y-%m-%dT%H:%MZ). Reset by deleting this file and re-running uci-defaults.)"
} > "$PASS_FILE"
chmod 600 "$PASS_FILE"

# --- Patch netdata.conf: bind localhost only + disable cloud/registry/anonymous stats ---
mkdir -p /etc/netdata
cat > "$NETDATA_CONF" <<'EOF'
[global]
    run as user = netdata
    memory mode = ram
    history = 3600
    update every = 2
    anonymous statistics = no

[web]
    bind to = 127.0.0.1
    allow connections from = localhost
    allow dashboard from = localhost
    allow badges from = *
    allow streaming from = none
    allow netdata.conf from = localhost
    allow management from = none
    enable gzip compression = yes

[registry]
    enabled = no

[cloud]
    enabled = no

[plugins]
    cgroups = no
    apps = no
    proc = yes
    diskspace = yes
    tc = no
    idlejitter = no
    python.d = no
    charts.d = no
    node.d = no
    go.d = no
EOF

# --- nginx reverse proxy with Basic Auth, served at /netdata/ ---
cat > "$NGINX_LOC" <<'EOF'
location /netdata/ {
    auth_basic           "Netdata - admin";
    auth_basic_user_file /etc/nginx/.netdata.htpasswd;
    proxy_set_header     Host $host;
    proxy_set_header     X-Real-IP $remote_addr;
    proxy_set_header     X-Forwarded-For $proxy_add_x_forwarded_for;
    proxy_pass           http://127.0.0.1:19999/;
    proxy_http_version   1.1;
    proxy_set_header     Connection "";
}
EOF

# --- Make sure firewall blocks 19999 from outside (defense in depth) ---
# (banip / fw4 default already drops; explicit rule below for clarity)
uci -q batch <<EOF
delete firewall.netdata_block 2>/dev/null || true
add firewall rule
set firewall.@rule[-1].name='Block-Netdata-Direct'
set firewall.@rule[-1].src='*'
set firewall.@rule[-1].proto='tcp'
set firewall.@rule[-1].dest_port='19999'
set firewall.@rule[-1].target='REJECT'
commit firewall
EOF

# Enable netdata at boot
/etc/init.d/netdata enable 2>/dev/null || true

logger -t netdata-harden "netdata bound to 127.0.0.1, basic-auth enabled, password in /root/.netdata-password"
exit 0
