# files/mtk-filogic-enterprise/etc/uci-defaults/93-services-disable
#!/bin/sh
# Disable services that require user-supplied configuration before they are safe to run.
# User explicitly enables them in LuCI after configuring them.

set -e

for svc in \
    mwan3 \
    frpc \
    homeproxy \
    sing-box \
    radvd \
    radius \
    freeradius3 \
    ; do
    if [ -x /etc/init.d/$svc ]; then
        /etc/init.d/$svc disable 2>/dev/null || true
        /etc/init.d/$svc stop 2>/dev/null || true
    fi
done

logger -t services-disable "disabled: mwan3 frpc homeproxy sing-box freeradius3 (user must enable after config)"
exit 0
